# launchd/com.remimarents.trening-reminder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.remimarents.trening-reminder</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/remimarents/Documents/Codex/2026-06-22/hjelp-meg-lage-en-webapp-for/scripts/trening-reminder.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/remimarents/Documents/Codex/2026-06-22/hjelp-meg-lage-en-webapp-for</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>APP_URL</key>
    <string>https://marents.no/trening/</string>
    <key>TZ_NAME</key>
    <string>Europe/Oslo</string>
    <key>DEFAULT_REMINDER_TIME</key>
    <string>19:30</string>
  </dict>

  <key>StartInterval</key>
  <integer>900</integer>

  <key>StandardOutPath</key>
  <string>/Users/remimarents/Library/Logs/trening-reminder.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/remimarents/Library/Logs/trening-reminder.err.log</string>
</dict>
</plist>
